# dist/systemd/bufu_filebroker_crash_notify.service
[Unit]
Description=Crash Noftifier for BUFU FileName Broker Server Service
Wants=network-online.target nfs.service
After=hltd.service network.service nfs.service rsyslog.service


[Service]
Type=oneshot
#User=notspecified
WorkingDirectory=/opt/bufu_filebroker/scripts
ExecStart=/opt/bufu_filebroker/scripts/crash_notify.sh
