# dist/systemd/bufu_filebroker.service
[Unit]
Description=BUFU FileName Broker Server Service
Wants=network-online.target
After=hltd.service network.service nfs.service rsyslog.service
OnFailure=bufu_filebroker_crash_notify.service

[Service]
Type=forking
#User=notspecified
WorkingDirectory=/cmsnfsscratch/globalscratch/pzejdl/bufu_filebroker
ExecStart=/cmsnfsscratch/globalscratch/pzejdl/bufu_filebroker/bufu_filebroker.sh
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target
